# debian/system76-driver.upstart
# system76-daemon - Support for Airplane Mode hotkey on System76 Laptops
#
# 

description	"System76 airplane-mode hotkey and LED support"

start on (login-session-start or desktop-session-start)
stop on desktop-shutdown

respawn

exec /usr/lib/system76-driver/system76-daemon

